# storage/.config/system.d/in.service
[Unit]
Description = GPIO monitor

[Service]
Type = idle
ExecStart = /storage/downloads/python/in.py

[Install]
WantedBy = multi-user.target
